# src/utils/SDNRDB-Export-Import-Scripts/crontab
*/2 * * * *  /home/nms-auto/Mohan/BackupScripts/export_db.sh /home/nms-auto/Mohan/BackupFolder
